# scripts/upgrade-tailscale/tailscaled.service
[Unit]
Description=Tailscale node agent
Documentation=https://tailscale.com/kb/
After=network-pre.target NetworkManager.service systemd-resolved.service
Wants=network-pre.target

[Service]
Type=notify
ExecStart=/usr/sbin/tailscaled --state=/var/lib/tailscale/tailscaled.state
Restart=on-failure
RuntimeDirectory=tailscale
RuntimeDirectoryMode=0755
StateDirectory=tailscale
StateDirectoryMode=0700

[Install]
WantedBy=multi-user.target
